# homebackup.service
[Unit]
Description=Backup service for local home directory

[Service]
Type=simple
ExecStart=/opt/geosscripts/rotatehomedirectory.sh

[Install]
WantedBy=multi-user.target